# vino-temp-control.service
[Unit]
Description=Vino Temperature Control System v0.5
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/vino_temp_control_v0.5
ExecStart=/usr/bin/python3 /home/pi/vino_temp_control_v0.5/app.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
